# backup/.config#autostart/able-live.desktop
[Desktop Entry]
Type=Application
Name=Able-Live
Exec=/usr/local/bin/able-live.sh
